# mecharm.service
[Unit]
Description=MechArm Control Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 -m server.app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
